# autoscale-source-split/01-monitoring-layer/systemd/autoscale-api.service
[Unit]
Description=Pre6G AutoScale API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=icclz2
Group=icclz2

WorkingDirectory=/home/icclz2/Pre6G/autoscale-source-split/03-shared-api-dashboard/autoscale_api

Environment="HOME=/home/icclz2"
Environment="KUBECONFIG=/home/icclz2/.kube/config"
Environment="PYTHONPATH=/home/icclz2/Pre6G:/home/icclz2/Pre6G/autoscale-source-split/03-shared-api-dashboard/autoscale_api"
EnvironmentFile=-/home/icclz2/Pre6G/autoscale-source-split/01-monitoring-layer/systemd/autoscale-api.env

ExecStart=/bin/bash -lc '/home/icclz2/Pre6G/iccl/bin/uvicorn app.main:app --host "${AUTOSCALE_API_HOST:-0.0.0.0}" --port "${AUTOSCALE_API_PORT:-8000}"'

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
